# deploy/ollama.service.d/override.conf
# ------------------------------------------------------------
# GCRM — Ollama systemd performance drop-in
# RAiTHE INDUSTRIES INCORPORATED © 2026
# ------------------------------------------------------------
#
# Tunes the local Ollama service so GCRM's concurrent LLM worker pool actually
# runs in parallel and the model stays hot. Without OLLAMA_NUM_PARALLEL, Ollama
# serializes requests internally and the GCRM worker pool gains nothing.
#
# INSTALL:
#   sudo mkdir -p /etc/systemd/system/ollama.service.d
#   sudo cp deploy/ollama.service.d/override.conf \
#           /etc/systemd/system/ollama.service.d/override.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart ollama
#
# VERIFY:
#   systemctl show ollama -p Environment
#   journalctl -u ollama -n 20 --no-pager
#
# TUNING:
#   • OLLAMA_NUM_PARALLEL must equal llm.concurrency in settings.yml (both 8).
#     Raise both together if the GPU has headroom; lower both if VRAM is tight.
#   • OLLAMA_MAX_LOADED_MODELS=2 lets the chat model and the embeddings model
#     (nomic-embed-text) stay resident together — only needed if you enable
#     llm.semantic_dedup. Leave at 2; it is harmless when dedup is off.

[Service]
# Serve up to 8 requests concurrently — match settings.yml llm.concurrency.
Environment="OLLAMA_NUM_PARALLEL=8"
# Keep the model resident indefinitely — no cold reloads between requests.
Environment="OLLAMA_KEEP_ALIVE=-1"
# Faster attention kernels on supported GPUs.
Environment="OLLAMA_FLASH_ATTENTION=1"
# Allow the chat model + embeddings model loaded together (for semantic_dedup).
Environment="OLLAMA_MAX_LOADED_MODELS=2"
