# deploy/ollama.service.d/override.conf
# ------------------------------------------------------------
# GCRM — Ollama systemd performance drop-in
# RAiTHE INDUSTRIES INCORPORATED © 2026
# ------------------------------------------------------------
#
# Tunes the local Ollama service so GCRM's concurrent LLM worker pool actually
# runs in parallel and the model stays hot. Without OLLAMA_NUM_PARALLEL, Ollama
# serializes requests internally and the GCRM worker pool gains nothing.
#
# INSTALL:
#   sudo mkdir -p /etc/systemd/system/ollama.service.d
#   sudo cp deploy/ollama.service.d/override.conf \
#           /etc/systemd/system/ollama.service.d/override.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart ollama
#
# VERIFY:
#   systemctl show ollama -p Environment
#   journalctl -u ollama -n 20 --no-pager
#
# TUNING:
#   • OLLAMA_NUM_PARALLEL must equal llm.concurrency in settings.yml (both 2 on the
#     8GB GTX 1080 — see below). Raise both together ONLY if the GPU has headroom;
#     lower both if VRAM is tight (a model split CPU/GPU in `ollama ps` = too high).
#   • OLLAMA_MAX_LOADED_MODELS=2 lets the chat model and the embeddings model
#     (nomic-embed-text) stay resident together — only needed if you enable
#     llm.semantic_dedup. Leave at 2; it is harmless when dedup is off.

[Service]
# Parallel requests — MUST match settings.yml llm.concurrency, and must fit VRAM:
# each slot adds a KV-cache copy. On this box (GTX 1080, 8GB) NUM_PARALLEL=8 pushed
# qwen2.5:7b to 8.7GB → 17% spilled to CPU → classify latency doubled (5.4s→11.3s,
# at the 12s timeout). 2 keeps qwen + nomic-embed-text fully on GPU (~6.1GB). Raise
# ONLY if `ollama ps` still shows 100% GPU and there's free VRAM headroom.
Environment="OLLAMA_NUM_PARALLEL=2"
# Keep the model resident indefinitely — no cold reloads between requests.
Environment="OLLAMA_KEEP_ALIVE=-1"
# Faster attention kernels on supported GPUs.
Environment="OLLAMA_FLASH_ATTENTION=1"
# Allow the chat model + embeddings model loaded together (for semantic_dedup).
Environment="OLLAMA_MAX_LOADED_MODELS=2"
